# sync-server/tplanner-sync.service
[Unit]
Description=tPlanner Sync Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/hamhuo/Documents/sync-server
ExecStart=/home/hamhuo/.nvm/versions/node/v22.22.3/bin/node /home/hamhuo/Documents/sync-server/server.js
Restart=on-failure
RestartSec=5

Environment=PORT=37401
Environment=DATA_DIR=/home/hamhuo/Documents/sync-server/data

# 以文件所有者身份运行，避免权限问题
User=hamhuo
Group=hamhuo

# 同时写 journal 和文件（append: 前缀是必须的）
StandardOutput=append:/home/hamhuo/Documents/sync-server/data/server.log
StandardError=append:/home/hamhuo/Documents/sync-server/data/server.log
SyslogIdentifier=tplanner-sync

[Install]
WantedBy=multi-user.target
